# rootfs/etc/services.d/apiserver/run
#!/bin/sh
export EPAPER_DIR=/srv/epaper
export EPAPER_ACCEL_DIR=/epaper/
exec apiserver
