# rootfs/etc/services.d/apiserver/run
#!/bin/sh
export EPAPER_DIR=/srv/epaper
exec apiserver
